# systemd/nvchecker.timer
[Unit]
Description=Run nvchecker every hour

[Timer]
OnActiveSec=1min
OnCalendar=hourly

[Install]
WantedBy=timers.target
